# apps/mucks/mucks.service
[Unit]
Description=mucks
After=network.target

[Service]
Type=simple
ExecStart=mucks
Restart=on-failure

[Install]
WantedBy=multi-user.target
